# templates/default/etc/systemd/system/weave.service.erb
[Unit]
Description=Weave Network
Documentation=http://docs.weave.works/weave/latest_release/
Requires=docker.service
After=docker.service
[Service]
Environment="WEAVE_NO_FASTDP=1"
ExecStartPre=/usr/local/bin/weave launch-router <%= @parameters %>
ExecStart=/usr/bin/docker attach weave
ExecStop=/usr/local/bin/weave stop-router
[Install]
WantedBy=multi-user.target